# cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Project Dashboard
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/home/opc/project-dashboard
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:3000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
